# debian/cloudagent.service
[Unit]
Description=CloudAgent
After=cloud-final.service

[Service]
ExecStart=/usr/bin/cloudagent -f

[Install]
WantedBy=multi-user.target
